# venom_network/venom-network.service
[Unit]
Description=Venom Network Daemon - WiFi/Bluetooth/Ethernet Service
Documentation=https://github.com/vaxp/venom-network
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/venom_network
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
